# firmware/Netgear_Centria_WNDR4700/squashfs-root/etc/init.d/bridge
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=15

start() {
	brctl addbr $BR_IF
	brctl setfd $BR_IF 0
	brctl stp $BR_IF 0
}
