# firmware/Netgear_Centria_WNDR4700/squashfs-root/etc/init.d/hotplug2
#!/bin/sh /etc/rc.common
START=70

boot() {
	insmod /lib/ufsd/ufsd.ko
	mkdir /tmp/mnt
#	mount -t sysfs sysfs /sys
	hotplug2 --persistent --coldplug &
}
