# firmware/Netgear_Centria_WNDR4700/squashfs-root/etc/init.d/madwifi
#!/bin/sh /etc/rc.common
START=75
boot() {
	echo 3 > /proc/sys/vm/drop_caches
#	/sbin/wlan up
}
