# firmware/Netgear_Centria_WNDR4700/squashfs-root/etc/init.d/pot
#!/bin/sh /etc/rc.common
START=65


pot_start() {
	potd
	potval
}

boot() {
	pot_start
}

pot_stop(){
	killall potval
	killall potd
}

restart(){
	pot_stop
	pot_start
}
