# firmware/Netgear_Centria_WNDR4700/squashfs-root/etc/init.d/traffic_meter.init
#!/bin/sh /etc/rc.common
START=85

boot() {
	/sbin/cmd_traffic_meter start
}
