# firmware/Netgear_Centria_WNDR4700/squashfs-root/etc/init.d/utelnetd
#!/bin/sh /etc/rc.common
START=60
boot() {
	# Start it by `telnetenable` in Router Debugging Mode.
	#/usr/sbin/utelnetd -d
	echo "Start utelnetd by telnetenable"
}
