# firmware/Netgear_Centria_WNDR4700/squashfs-root/etc/crontabs/broker
*/1 * * * * /etc/init.d/brokerd.sh start
