# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/ap-mode.list
/etc/init.d/ap-mode
/usr/sbin/restart_ap_udhcpc
/usr/share/zcip/zcip.script
