# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/atheros-aquila-madwifi-11n-scripts.list
/etc/ath.orig/apdown
/etc/ath
/etc/init.d/madwifi-11n-scripts
/etc/ath.orig/WSC_sta.conf
/sbin/update-wifi
/usr/lib/stainfo.awk
/etc/ath.orig/hostapd0.7.0_conf/WSC.conf
/etc/ath.orig/wifi_schedule
/etc/ath.orig/hostapd0.7.0_conf/PSK.ap_bss
/etc/dni-wifi-config
/etc/ath.orig/board.conf
/etc/ath.orig/wifi_toggle
/etc/ath.orig/apcfg
/etc/ath.orig/WSC_WEP.conf
/etc/ath.orig/PSK.ap_bss
/etc/ath.orig/wpscli_sta
/etc/ath.orig/PSK.sta
/sbin/wlan_old
/etc/ath.orig/activateVAP
/etc/ath.orig/killVAP
/etc/ath.orig/wpscli_ap
/etc/rc.d/rc.wlan
/sbin/dni_apup
/etc/ath.orig/apup
/etc/ath.orig/hostapd0.7.0_conf/EAP.ap_bss
/etc/ath.orig/EAP.ap_bss
/lib/wifi/atheros11n.sh
/etc/ath.orig/RADIUS_MAC.ap_bss
/etc/ath.orig/WSC.conf
/etc/ath.orig/makeVAP
