# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/base-files.list
/bin/login
/etc/resolv.conf
/etc/init.d/net-wan
/usr/sbin/check_smart_error.sh
/etc/init.d/done
/var
/etc/nsswitch.conf
/etc/hotplug2-init.rules
/lib/firstboot/99_10_no_fo_cleanup
/etc/preinit
/lib/firstboot/50_pivot
/etc/init.d/network
/sbin/mount_root
/sbin/backup
/usr/sbin/format_sata
/lib/firstboot/30_reset_copy_rom
/etc/hotplug2-common.rules
/etc/shells
/lib/preinit/70_initramfs_test
/lib/preinit/50_choose_console
/lib/firstboot/10_no_fo_clear_overlay
/etc/protocols
/lib/preinit/20_check_jffs2_ready
/sbin/ip_mac
/lib/preinit/10_indicate_preinit
/etc/init.d/S41App_start.sh
/sbin/detect_factory
/etc/hotplug.d/iface/00-netstate
/lib/firstboot/99_10_with_fo_cleanup
/sbin/ifdown
/etc/config/system
/usr/share/udhcpc/default.script
/lib/network/config.sh
/lib/preinit/30_device_fs_daemons
/etc/sysctl.conf
/lib/preinit/90_mount_no_jffs2
/sbin/sysupgrade
/lib/preinit/70_pivot_jffs2_root
/lib/preinit/99_10_run_init
/lib/firstboot/20_has_mini_fo
/lib/firstboot/10_reset_has_mini_fo
/lib/firstboot/10_determine_parts
/sbin/printhosts
/sbin/led.sh
/etc/hotplug.d/net/10-net
/lib/preinit/50_indicate_regular_preinit
/lib/firstboot/40_no_fo_copy_ramoverlay
/etc/banner
/etc/init.d/net-lan
/lib/preinit/10_essential_fs
/etc/config/timeserver
/etc/rc.common
/lib/preinit/10_indicate_failsafe
/lib/preinit/40_init_shm
/init
/etc/init.d/urandom
/etc/init.d/net-iptv
/lib/preinit/40_mount_jffs2
/lib/preinit/20_device_fs_mount
/lib/preinit/80_mount_root
/etc/rc.local
/etc/init.d/rc.firewall-iptables
/etc/init.d/mountSATA
/etc/init.d/sysctl
/lib/firstboot/30_is_rootfs_mounted
/sbin/hotplug-call
/sbin/log_fan_rpm_fault
/lib/ufsd/ufsd.ko
/sbin/usb_disk_event
/lib/preinit/40_mount_devpts
/lib/preinit/90_init_console
/lib/preinit/99_10_mount_no_mtd
/lib/preinit/10_check_for_mtd
/lib/upgrade/common.sh
/sbin/printnss
/etc/sysupgrade.conf
/lib/firstboot/05_firstboot_skip
/sbin/cmd_autoreset
/etc/init.d/bridge
/lib/upgrade/keep.d/base-files
/etc/preinit.arch
/lib/preinit/40_run_failsafe_hook
/etc/group
/lib/preinit/30_failsafe_wait
/etc/functions.sh
/sbin/firstboot
/etc/hosts
/etc/hotplug.d/usb/10-usb
/etc/init.d/upnp
/lib/preinit/60_init_hotplug
/etc/config/network
/etc/init.d/rcS
/sbin/wps_pbc
/usr/lib/common.awk
/sbin/ifup
/rom/note
/etc/inittab
/sbin/qos.sh
/lib/preinit/05_mount_skip
/sbin/wifi
/etc/init.d/switch
/sbin/run-ramfs
/sbin/cmd_ddns
/etc/profile
/etc/hotplug.d/iface/10-routes
/lib/preinit/99_10_failsafe_login
/etc/hotplug.d/iface/40-rdate
/etc/TZ
/lib/firstboot/20_no_fo_mount_jffs
/lib/preinit/90_restore_config
/lib/upgrade/keep.d/base-files-essential
/lib/firstboot/30_no_fo_pivot
/etc/passwd
/etc/hotplug.d/ieee1394/10-ieee1394
/etc/mtab
/mnt
/lib/firstboot/40_copy_ramoverlay
/usr/sbin/luns_scan.sh
/etc/init.d/boot
/sbin/check_HDD_capacity
/lib/functions/boot.sh
/bin/ipcalc.sh
/etc/diag.sh
/sbin/storage_test
/etc/init.d/rc.wlan
/lib/firstboot/20_reset_clear_jffs
/lib/preinit/41_merge_overlay_hooks
