# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/dnsmasq.list
/etc/dnsmasq.conf
/etc/init.d/dnsmasq
/usr/sbin/dns-hijack
/usr/sbin/dnsmasq
