# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/hotplug2.list
/sbin/app_mount
/sbin/umount_unapproved
/sbin/hotplug2
/etc/init.d/hotplug2
/etc/hotplug2.rules
/sbin/hotplug2.umount
/sbin/hotplug2.mount
/sbin/hotplug2.testbusy
