# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/igmpproxy.list
/sbin/cmdigmp
/etc/init.d/igmpproxy.init
/sbin/igmpproxy
