# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/kmod-kcode.list
/lib/modules/2.6.32.11-wndr4700/NetUSB.ko
/etc/init.d/kcode
/lib/modules/2.6.32.11-wndr4700/GPL_NetUSB.ko
