# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/kmod-wndr4700-madwifi.list
/sbin/pktlogconf
/lib/modules/2.6.32.11-wndr4700/ath_dfs.ko
/sbin/dumpregs
/lib/modules/2.6.32.11-wndr4700/asf.ko
/lib/modules/2.6.32.11-wndr4700/ath_pktlog.ko
/sbin/athstats
/lib/modules/2.6.32.11-wndr4700/ath_dev.ko
/lib/modules/2.6.32.11-wndr4700/adf.ko
/etc/init.d/madwifi
/sbin/pktlogdump
/lib/modules/2.6.32.11-wndr4700/ath_hal.ko
/sbin/80211stats
/lib/modules/2.6.32.11-wndr4700/umac.ko
/lib/modules/2.6.32.11-wndr4700/wlan_ext.ko
/sbin/reg
/lib/modules/2.6.32.11-wndr4700/ath_rate_atheros.ko
/sbin/cfg
/sbin/athstatsclr
/sbin/cgiMain
/sbin/radartool
/sbin/wlanconfig
