# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/netatalk.list
/usr/config/netatalk/AppleVolumes.default
/usr/lib/uams/uams_dhx.so
/etc/netatalk/AppleVolumes.system
/usr/sbin/afppasswd
/etc/init.d/afpd
/usr/lib/uams/uams_passwd.so
/etc/netatalk/afppasswd
/etc/netatalk/AppleVolumes.default
/usr/sbin/cnid_metad
/usr/lib/uams/uams_randnum.so
/usr/lib/uams/uams_dhx2_passwd.so
/usr/sbin/detach_afp_shares
/etc/netatalk/afp_signature.conf
/usr/lib/uams/uams_guest.so
/usr/lib/uams/uams_dhx_passwd.so
/usr/sbin/cnid_dbd
/usr/lib/uams/uams_dhx2.so
/etc/netatalk/afpd.conf
/usr/sbin/afpd
