# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/radvd.list
/usr/sbin/radvdump
/etc/radvd.conf
/usr/sbin/radvd
/etc/init.d/radvd
/etc/config/radvd
