# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/readyshare.list
/opt/remote/bin/register.sh
/opt/remote/bin/comm.sh~
/etc/crontabs/broker
/opt/remote/bin/p2p.sh
/etc/init.d/leafp2p.sh
/opt/remote/bin/env.sh
/opt/remote/bin/comm.sh
/opt/remote/bin/hook.sh
/opt/remote/bin/unregister.sh
/opt/remote/bin/leafp2p
/opt/remote/bin/fscmd.sh~
/opt/remote/bin/syscmd.sh~
/opt/remote/bin/checkleafnets.sh~
/etc/init.d/remote.sh
/opt/remote/bin/syscmd.sh
/opt/remote/bin/rules.sh
/opt/remote/bin/checkleafp2p.sh
/opt/remote/bin/tail
/etc/init.d/brokerd.sh
/opt/remote/bin/checkleafnets.sh
/opt/remote/bin/fscmd.sh
/opt/remote/bin/enclosure.sh
