# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/uhttpd.list
/usr/sbin/uhttpd
/www/cgi-bin/uhttpd.sh
/etc/init.d/uhttpd
/etc/config/uhttpd
/lib/upgrade/keep.d/uhttpd
