# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/utelnetd.list
/etc/init.d/utelnetd
/usr/sbin/utelnetd
