# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/wlan-common.list
/etc/hotplug.d/wps/00-wps
/etc/hotplug.d/wps/01-wps-led
/sbin/wpsled
/etc/init.d/wlan-common
/etc/config/wireless
/sbin/wlan
