# firmware/Netgear_Centria_WNDR4700/squashfs-root/usr/lib/opkg/info/zebra.list
/usr/sbin/ripngd
/usr/sbin/zebra
/etc/init.d/ripngd
/etc/ripngd.conf
/etc/init.d/zebra
/etc/zebra.conf
